# package/plugin-gargoyle-ap-management/files/etc/uci-defaults/init_gargoyle_ap_management_plugin
#!/bin/sh

uci set gargoyle.display.connection_ap_management="AP Management"
uci set gargoyle.scripts.connection_ap_management="ap_management.sh"
uci set gargoyle.connection.ap_management="700"

uci commit
